# nexustrader-mcp-sse.service
# =============================================================================
# 可选：systemd user service 模板（仅供高级用户手动配置，非必须）
#
# 此文件不会被自动安装，OpenClaw Skill 也不依赖它。
# 如需使用 systemd 管理服务器（开机自启），请按以下步骤手动操作：
#
#   1. 将 %h 替换为您的 $HOME 路径，然后复制到 systemd 目录：
#      sed "s|%h|$HOME|g" nexustrader-mcp-sse.service \
#          > ~/.config/systemd/user/nexustrader-mcp-sse.service
#
#   2. 启用并启动：
#      systemctl --user daemon-reload
#      systemctl --user enable --now nexustrader-mcp-sse
#
#   3. 查看状态：
#      systemctl --user status nexustrader-mcp-sse
#
# 不使用 systemd 的用户：直接运行 `uv run nexustrader-mcp start` 即可。
# =============================================================================

[Unit]
Description=NexusTrader MCP Server (SSE mode, port 18765)
Documentation=https://github.com/nexustrader/nexustrader-mcp
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
ExecStart=uv \
    --directory %h/NexusTrader-mcp \
    run --python 3.11 \
    nexustrader-mcp serve \
    --host 127.0.0.1 \
    --port 18765 \
    --config %h/NexusTrader-mcp/config.yaml
Restart=on-failure
RestartSec=15
# Give the trading engine enough time to initialize before systemd declares it failed
TimeoutStartSec=120
# Tighten environment to avoid conda/venv interference
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment=PYTHONPATH=
Environment=PYTHONHOME=
Environment=CONDA_PREFIX=
Environment=CONDA_DEFAULT_ENV=
Environment=CONDA_SHLVL=0
Environment=UV_PYTHON_PREFERENCE=only-managed
Environment=UV_PYTHON=cpython-3.11
# Redirect stdout/stderr to the journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexustrader-mcp-sse

[Install]
WantedBy=default.target
